# cyrus-imapd.cron
#!/bin/sh
su cyrus -s /bin/sh -c '/usr/lib/cyrus/deliver -E 3'
